# ansible/files/certbot.service
[Unit]
Description=Certbot Renew Service

[Service]
Type=oneshot
ExecStart=docker exec production_nginx certbot renew

[Install]
WantedBy=multi-user.target